# src/main/assembly/dist/bin/easy-test.service
[Unit]
Description=EASY Test Service

[Service]
ExecStart=/opt/dans.knaw.nl/easy-test/bin/easy-test-service-start

User=easy-test
Group=easy-test

[Install]
WantedBy=multi-user.target
